# docs/code_examples/supervisord.conf
[program:errbot]
directory=/etc/errbot
command=/path/to/my/virtualenv/bin/errbot
autostart=true
autorestart=true
